# share/systemd/system/gate-runtime.socket
[Unit]
Description=Gate runtime socket

[Socket]
ListenStream=/run/gate/runtime.sock

[Install]
WantedBy=multi-user.target
